# khoji/com.lekhio.khoji.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Khoji launchd job. Runs the watcher once a day at 05:15, isolated under its own
  label so it never clashes with the personal bot already on the mini.

  The paths below point at the folder that CONTAINS run.sh. That is
  /Users/jagchahil/lekhio-khoji. Between 7 and 12 July every path here carried an
  extra /khoji on the end, pointing at an empty subfolder, so the daily job fired
  into a void and never ran once. If you edit these, check run.sh actually exists
  at the path you have written.

  Install:
    cp com.lekhio.khoji.plist ~/Library/LaunchAgents/
    launchctl unload ~/Library/LaunchAgents/com.lekhio.khoji.plist 2>/dev/null
    launchctl load ~/Library/LaunchAgents/com.lekhio.khoji.plist
  Run once now:   launchctl kickstart -k gui/$(id -u)/com.lekhio.khoji
  Stop / remove:  launchctl unload ~/Library/LaunchAgents/com.lekhio.khoji.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lekhio.khoji</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jagchahil/lekhio-khoji/run.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jagchahil/lekhio-khoji</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>5</integer>
    <key>Minute</key><integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jagchahil/lekhio-khoji/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jagchahil/lekhio-khoji/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
